# recon-fs/var/lib/dpkg/info/linux-rock64-package.list
/.
/boot
/boot/efi
/boot/efi/extlinux
/boot/efi/extlinux/extlinux.conf
/lib
/lib/firmware
/lib/firmware/rtlbt
/lib/firmware/rtlbt/rtl8723a_config
/lib/firmware/rtlbt/rtl8723b_fw
/lib/firmware/rtlbt/rtl8703b_config
/lib/firmware/rtlbt/rtl8723b_config
/lib/firmware/rtlbt/rtl8723a_fw
/lib/firmware/rtlbt/rtl8703b_fw
/lib/firmware/rtlbt/rtl8723b_VQ0_fw
/lib/firmware/rtlbt/rtl8723cs_cg_config
/lib/firmware/rtlbt/rtl8723cs_xx_config
/lib/firmware/rtlbt/rtl8723cs_xx_fw
/lib/firmware/rtlbt/rtl8703a_fw
/lib/firmware/rtlbt/rtl8723b_config_2Ant_S0
/lib/firmware/rtlbt/rtl8703a_config
/lib/firmware/rtlbt/rtl8723b_VQ0_config
/lib/firmware/rtlbt/rtl8723cs_vf_fw
/lib/firmware/rtlbt/rtl8723cs_vf_config
/lib/firmware/rtlbt/rtl8723cs_cg_fw
/etc
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/50-mail400.rules
/etc/udev/rules.d/50-hevc.rules
/etc/udev/rules.d/hdmi.rules
/etc/udev/rules.d/50-hevc-rk3399.rules
/etc/udev/rules.d/50-vpu-rk3399.rules
/etc/udev/rules.d/50-vpu.rules
/etc/udev/rules.d/50-mail.rules
/etc/chromium-browser
/etc/chromium-browser/customizations
/etc/chromium-browser/customizations/rock64-optimisations
/etc/modprobe.d
/etc/modprobe.d/.gitkeep
/etc/update-motd.d
/etc/update-motd.d/05-figlet
/etc/lightdm
/etc/lightdm/lightdm-gtk-greeter.conf.d
/etc/lightdm/lightdm-gtk-greeter.conf.d/99_xxx_rock64.conf
/etc/systemd
/etc/systemd/system
/etc/systemd/system/restore-sound-after-resume.service
/etc/systemd/system/first-boot.service
/etc/systemd/system/ssh-keygen.service
/etc/systemd/system/rtk-hciattach.service
/etc/systemd/system/restart-network-manager-after-resume.service
/etc/systemd/system/store-sound-on-suspend.service
/etc/cron.d
/etc/cron.d/make_nas_processes_faster
/etc/cron.d/startup
/etc/network
/etc/network/interfaces.d
/etc/network/interfaces.d/eth0
/etc/network/if-up.d
/etc/network/if-up.d/rock64-offload
/etc/X11
/etc/X11/xorg.conf.d
/etc/X11/xorg.conf.d/20-modesetting.conf
/etc/firmware
/etc/firmware/BCM4335C0.hcd
/etc/firmware/BCM4343A0 26M.hcd
/etc/firmware/4343A0.hcd
/etc/modules-load.d
/etc/modules-load.d/.gitkeep
/usr
/usr/local
/usr/local/bin
/usr/local/bin/hdmi-toggle
/usr/local/sbin
/usr/local/sbin/rock64_first_boot.sh
/usr/local/sbin/rock64_eth0_stats.sh
/usr/local/sbin/rock64_health.sh
/usr/local/sbin/install_openmediavault.sh
/usr/local/sbin/install_desktop.sh
/usr/local/sbin/rock64_fix_performance.sh
/usr/local/sbin/rock64_diagnostics.sh
/usr/local/sbin/resize_rootfs.sh
/usr/share
/usr/share/flash-kernel
/usr/share/flash-kernel/db
/usr/share/flash-kernel/db/rock64.db
/usr/share/alsa
/usr/share/alsa/init
/usr/share/alsa/init/rt5640
/usr/share/alsa/init/rt5651
/usr/share/alsa/init/10rockchip
/usr/share/alsa/init/rt5616
/usr/share/backgrounds
/usr/share/backgrounds/ubuntu-mate-rock64
/usr/share/backgrounds/ubuntu-mate-rock64/ROCK64-Wallpaper-1.jpg
/usr/share/backgrounds/ubuntu-mate-rock64/ROCK64-Wallpaper-6.jpg
/usr/share/backgrounds/ubuntu-mate-rock64/ROCK64-Wallpaper-2.jpg
/usr/share/backgrounds/ubuntu-mate-rock64/ROCK64-Wallpaper-4.jpg
/usr/share/backgrounds/ubuntu-mate-rock64/ROCK64-Wallpaper-3.jpg
/usr/share/backgrounds/ubuntu-mate-rock64/ROCK64-Wallpaper-5.jpg
/usr/share/doc
/usr/share/doc/linux-rock64-package
/usr/share/doc/linux-rock64-package/changelog.gz
/usr/local/sbin/armbianmonitor
