# recon-fs/etc/update-motd.d/05-figlet
#!/bin/sh

#figlet $(hostname)
figlet "Recon Sentinel"
